# pkgs/scripts/rpm/preremove
if [ $1 -eq 0 ]; then
  # Package removal, not upgrade
  /bin/systemctl --no-reload disable sftpgo.service >/dev/null 2>&1 || :
  /bin/systemctl stop sftpgo.service >/dev/null 2>&1 || :
fi
